# fan.service
[Unit]
Description=Start Fan

[Service]
Type=simple
ExecStart=/bin/sh /usr/local/bin/fan.sh

[Install]
WantedBy=multi-user.target